# samplicator_configs/systemd/system/samplicator-ipfix.service.d/samplicator-ipfix.conf
[Service]
# Samplicator - IPFIX
Environment="SAMP_IPFIX_LISTEN_IP=0.0.0.0"
Environment="SAMP_IPFIX_LISTEN_PORT=4739"
