# samplicator_configs/systemd/system/samplicator-snmptrap.service.d/samplicator-snmptrap.conf
[Service]
# Samplicator - SNMP Trap
Environment="SAMP_SNMPTRAP_LISTEN_IP=0.0.0.0"
Environment="SAMP_SNMPTRAP_LISTEN_PORT=162"
