# samplicator_config/systemd/system/samplicator-sflow.service.d/samplicator-sflow.conf
[Service]
# Samplicator - sFlow
Environment="SAMP_SFLOW_LISTEN_IP=0.0.0.0"
Environment="SAMP_SFLOW_LISTEN_PORT=6343"
