# samplicator_config/profile.d/samplicator-sflow.sh
# Samplicator - sFlow
export SAMP_SFLOW_LISTEN_IP=0.0.0.0
export SAMP_SFLOW_LISTEN_PORT=6343
