# samplicator_configs/profile.d/samplicator-ipfix.sh
# Samplicator - IPFIX
export SAMP_IPFIX_LISTEN_IP=0.0.0.0
export SAMP_IPFIX_LISTEN_PORT=4739
