# samplicator_config/systemd/system/samplicator-netflow.service.d/samplicator-netflow.conf
[Service]
# Samplicator - Netflow
Environment="SAMP_NETFLOW_LISTEN_IP=0.0.0.0"
Environment="SAMP_NETFLOW_LISTEN_PORT=2055"
