# samplicator_configs/systemd/system/samplicator-syslog.service.d/samplicator-syslog.conf
[Service]
# Samplicator - Syslog
Environment="SAMP_SYSLOG_LISTEN_IP=0.0.0.0"
Environment="SAMP_SYSLOG_LISTEN_PORT=514"
