# telegraf_config/profile.d/telegraf.sh
# Telegraf
export TELEGRAF_OUTPUT_SOCKET_ADDRESS=tcp://<LOGSTASH_IP>:58086
